# jee-duo.service.example
Copy this to /etc/systemd/system/jee-duo.service on your VPS, adjust the
paths/user, then run:

    sudo systemctl daemon-reload
    sudo systemctl enable --now jee-duo

--------------------------------------------------------------------------

[Unit]
Description=JEE Duo study server
After=network.target

[Service]
Type=simple
User=YOUR_LINUX_USERNAME
WorkingDirectory=/path/to/jee-duo-server
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
EnvironmentFile=/path/to/jee-duo-server/.env

[Install]
WantedBy=multi-user.target
